# scripts/team-briefing/com.opsconsole.team-briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.opsconsole.team-briefing</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/yss/개발/build/OPS-Console/scripts/team-briefing/publish-local.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/yss/개발/build/OPS-Console</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yss/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>5</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/yss/Library/Logs/team-briefing.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yss/Library/Logs/team-briefing.log</string>
</dict>
</plist>
